# ansible/roles/enlarge_rx_ring_buffer/templates/99-enlarge-rx-ring-buffer.rules.jinja2
{% for key in network_config %}
SUBSYSTEM=="net",ACTION=="add|change",KERNEL=="{{ network_config[key].interface_name }}",RUN+="/usr/sbin/ethtool -G %k rx {{ rx_ring_buffer_size }}"
{% endfor %}
